# .docker/app/crontab
* * * * * cd /var/www && php artisan schedule:run >> /dev/null 2>&1
*/30 * * * * cd /var/www && php artisan jobs:start >> /dev/null 2>&1